# etc/services.d/mai-inventory-mgr/finish
#!/usr/bin/execlineb -S0

# Stops all services when "mai-inventory-mgr" stops.
s6-svscanctl -t /var/run/s6/services