# etc/services.d/mai-inventory-mgr/run
#!/usr/bin/with-contenv bash
sleep 5 # Give litestream a chance to start
exec python3.7 /app/run.py