# etc/services.d/mai-inventory-mgr/run
#!/usr/bin/with-contenv bash
exec python3.7 /app/run.py